# conf/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:apache]
command=apache2-foreground

[program:cron]
command = /usr/sbin/cron -f -L 15
stdout_logfile  = /var/log/cron.log
stderr_logfile  = /var/log/cron.log
autorestart=true
